# deploy/cortex.service
[Unit]
Description=Cortex — Product Knowledge Graph API
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/cortex
ExecStart=/usr/bin/docker compose up
ExecStop=/usr/bin/docker compose down
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cortex

[Install]
WantedBy=multi-user.target
